# install/conf/opsany-paas/mcp-server/mcp-supervisor.ini
[program: mcp-server]
command = python3 /opt/opsany/opsany-mcp-server/server.py --config /opt/opsany/opsany-mcp-server/config/config.yaml
stdout_logfile = /opt/opsany/logs/opsany-mcp-server/server.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
